# docker/conf.d/dnsmasq.conf
[program:dnsmasq]
command=dnsmasq -C /etc/dnsmasq.conf -k -d
