# docker/conf.d/docker.conf
[program:dockerd]
command=/usr/bin/dockerd -H unix:///var/run/docker.sock -H tcp://127.0.0.1:2375
