# docker/conf.d/vpnkit.conf
[program:vpnkit]
command=/usr/local/bin/wsl-vpnkit
