# docker/root/etc/supervisor/conf.d/chrony.conf
[program:chrony]
priority=100
command=chronyd -d
stderr_logfile = /var/log/supervisord/chrony-stderr.log
stdout_logfile = /var/log/supervisord/chrony-stdout.log
