# docker/root/etc/supervisor/conf.d/chrony.conf
[program:chrony]
priority=100
command=chrony -d
