# docker/root/etc/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
priority=200
command=dnsmasq -C /etc/dnsmasq.conf -k -d
stderr_logfile = /var/log/supervisor/dnsmasq-stderr.log
stdout_logfile = /var/log/supervisor/dnsmasq-stdout.log
