# docker/root/etc/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
priority=200
command=dnsmasq -C /etc/dnsmasq.conf -k -d
