# docker/root/etc/supervisor/conf.d/mount.conf
[program:mount]
priority=5
command=/usr/local/bin/wsl-mount
stderr_logfile = /var/log/supervisord/mount-stderr.log
stdout_logfile = /var/log/supervisord/mount-stdout.log
