# docker/root/etc/supervisor/conf.d/mount.conf
[program:mount]
priority=5
command=/usr/local/bin/wsl-mount
