# docker/root/etc/supervisor/conf.d/startup.conf
[program:startup]
priority=1
startsecs=0
command=/usr/local/bin/startup
stderr_logfile = /var/log/supervisord/startup-stderr.log
stdout_logfile = /var/log/supervisord/startup-stdout.log
