# docker/root/etc/supervisor/conf.d/startup.conf
[program:startup]
priority=1
startsecs=0
command=/usr/local/bin/startup
