# docker/root/etc/supervisor/conf.d/startup.conf
[program:startup]
priority=1
command=/usr/local/bin/startup
