# docker/root/etc/supervisor/conf.d/vpnkit.conf
[program:vpnkit]
priority=10
command=/usr/local/bin/wsl-vpnkit
