# docker/root/etc/supervisor/conf.d/vpnkit.conf
[program:vpnkit]
priority=10
command=/usr/local/bin/wsl-vpnkit
stderr_logfile = /var/log/supervisor/vpnkit-stderr.log
stdout_logfile = /var/log/supervisor/vpnkit-stdout.log
