# docker/root/etc/systemd/system/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket
# Do not stop on shutdown
DefaultDependencies=no
Wants=sysinit.target
After=sysinit.target

[Socket]
ListenStream=/run/dbus/system_bus_socket
ListenStream=127.0.0.1:12434
